# resources/upstart/tweetarchive-index.conf
description "tweet-archive.me indexing job"
author      "Tolleiv Nietsch"

# mount have to be ready for us:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/node"
    exec sudo -u node /usr/local/bin/node /home/vagrant/app/tweet-archive/index.js >> /var/log/tweetarchive-index.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
